# sleepcycle.service
[Service]
ExecStart=/home/kandre/src/sleepcycle/listen-and-publish.py
Restart=always

[Install]
WantedBy=multi-user.target
